# centos-pg/setup-pg/config/pg.conf
[program:postgresql]
command=/run/pgstart
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/lib/pgsql/9.2/data/pg_log/pgstartup.log

